# containers/rootfs/etc/s6-overlay/s6-rc.d/tv/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
cd /app
s6-setuidgid nextjs
node start.js